# battery_notify/.config/systemd/user/battery_monitor.timer
[Timer]
OnUnitInactiveSec=3

[Install]
WantedBy=timers.target
